# service/sftp-downloader.timer
[Unit]
Description=Timer for SFTP downloader script
 
[Timer]
OnUnitActiveSec=30min
OnBoot=10s
 
[Install]
WantedBy=timers.target